# systemd/power-performance.target
[Unit]
Description=System power management: performance mode
DefaultDependencies=no
Conflicts=power-save.target
