# etc/systemd/system/power-performance.target
[Unit]
Description=System power management: performance mode
Conflicts=power-save.target
