# etc/systemd/system/power-save.target
[Unit]
Description=System power management: power-saving mode
Conflicts=power-performance.target
After=power-performance.target
