# etc/rules.d/50-powersave.rules
ACTION=="change", SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_ONLINE}=="[01]", RUN+="/usr/local/lib/powersave/bin/power %E{POWER_SUPPLY_ONLINE}"
